# systemd/vllm-proxy.service
[Unit]
Description=vLLM Studio Proxy Server
After=network.target vllm-studio.service

[Service]
Type=simple
User=ser
WorkingDirectory=/home/ser/lmvllm
Environment="PATH=/home/ser/lmvllm/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ser/lmvllm/.venv/bin/python -m uvicorn proxy.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
